# assets/services/radm-health.timer
[Unit]
Description=RADM Health Timer (5min)
[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target